# CMake/Package/postinst
systemctl daemon-reload
systemctl enable eLinkDisplayStatus.service
systemctl start eLinkDisplayStatus.service
