# CMake/Package/prerm
systemctl stop eLinkDisplayStatus.service
systemctl disable eLinkDisplayStatus.service
